# deploy/compliance-api.service
[Unit]
Description=152check API (FastAPI)
After=network.target

[Service]
WorkingDirectory=/var/www/152check
EnvironmentFile=/var/www/152check/.env
ExecStart=/var/www/152check/.venv/bin/uvicorn api.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=3
User=www-data

[Install]
WantedBy=multi-user.target
